# deploy/linux/oo-agent.service
[Unit]
Description=oo-agent metrics collector
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=simple
ExecStart=/opt/oo-agent/venv/bin/oo-agent
Restart=on-failure
RestartSec=10
# Sensors, SMART and GPU BAR reads need root; everything else is
# best-effort and degrades gracefully without it.
User=root
Nice=10
MemoryMax=256M

[Install]
WantedBy=multi-user.target
